# supervisord.conf
[supervisord]
nodaemon=true

[program:python]
command=python3 -u /usr/local/bin/report-generator.py
autostart=true
autorestart=true
redirect_stderr=true
